# entryagent/etc/systemd/system/entryagent.service
[Unit]
Description=EntryBox Heartbeat Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=/etc/entryagent.env
ExecStart=/usr/bin/python3 /opt/entryagent/agent.py
Restart=always
RestartSec=5
Nice=10

[Install]
WantedBy=multi-user.target
